# deployment/systemd/azl-tests.service
[Unit]
Description=AZL Test Runner (Strict Mode)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/abdulrahman-alzalameh/azl-language
Environment=AZL_STRICT=1
ExecStart=/usr/bin/bash -lc '/home/abdulrahman-alzalameh/azl-language/scripts/run_all_tests.sh'
User=abdulrahman-alzalameh
Group=abdulrahman-alzalameh
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target


