# deployment/systemd/azl-tests.timer
[Unit]
Description=Daily AZL Test Runner Timer

[Timer]
OnCalendar=daily
Persistent=true
Unit=azl-tests.service

[Install]
WantedBy=timers.target


