# deployment/systemd/azl-runtime.service
[Unit]
Description=AZL Runtime (sysproxy-backed HTTP loop)
After=network.target azl-sysproxy.service
Requires=azl-sysproxy.service

[Service]
Type=notify
NotifyAccess=all
User=abdulrahman-alzalameh
WorkingDirectory=/home/abdulrahman-alzalameh/azl-language
EnvironmentFile=-/etc/default/azl
ExecStartPre=/bin/bash -lc 'mkdir -p .azl'
ExecStart=/bin/bash scripts/start_azl_native_mode.sh
Restart=always
RestartSec=2s
WatchdogSec=30s
StandardOutput=append:/home/abdulrahman-alzalameh/azl-language/.azl/daemon.out
StandardError=append:/home/abdulrahman-alzalameh/azl-language/.azl/daemon.err
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectKernelTunables=yes
ProtectControlGroups=yes
ProtectClock=yes
LockPersonality=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
UMask=0077
LimitNOFILE=8192
TimeoutStopSec=10s
KillMode=process

[Install]
WantedBy=multi-user.target
