# deployment/systemd/azl-daemon.service
[Unit]
Description=AZL Headless Daemon (pure AZL, HTTP API)
After=network-online.target azl-sysproxy.service
Wants=network-online.target azl-sysproxy.service

[Service]
Type=simple
User=abdulrahman-alzalameh
WorkingDirectory=/home/abdulrahman-alzalameh/azl-language
EnvironmentFile=-/etc/default/azl
EnvironmentFile=-%h/.config/azl/env
# Hardened defaults; allow override via /etc/default/azl
Environment=AZL_BUILD_API_PORT=8080
Environment=AZL_LOCAL_MODE=0
Environment=AZL_REQUIRE_API_TOKEN=true
Environment=AZL_HTTP_RATELIMIT_WINDOW_MS=1000
Environment=AZL_HTTP_RATELIMIT_MAX=200
Environment=AZL_SNAPSHOT_STALE_SEC=10
Environment=AZL_LOG_LEVEL=info
Environment=TMPDIR=%h/azl-language/.azl/tmp

ExecStartPre=/usr/bin/bash -lc 'mkdir -p .azl .azl/tmp'
ExecStart=/usr/bin/bash -lc 'AZL_SYSTEMD=1 AZL_NODE_RUNTIME=0 AZL_NODE_HTTP=0 AZL_BUILD_API_PORT=${AZL_BUILD_API_PORT} AZL_LOCAL_MODE=${AZL_LOCAL_MODE} AZL_REQUIRE_API_TOKEN=${AZL_REQUIRE_API_TOKEN} AZL_HTTP_RATELIMIT_WINDOW_MS=${AZL_HTTP_RATELIMIT_WINDOW_MS} AZL_HTTP_RATELIMIT_MAX=${AZL_HTTP_RATELIMIT_MAX} AZL_SNAPSHOT_STALE_SEC=${AZL_SNAPSHOT_STALE_SEC} AZL_LOG_LEVEL=${AZL_LOG_LEVEL} bash scripts/run_enterprise_daemon.sh'

Restart=always
RestartSec=2s
TimeoutStopSec=15s
KillMode=control-group
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectKernelTunables=yes
ProtectControlGroups=yes
ProtectClock=yes
LockPersonality=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
UMask=0077
LimitNOFILE=8192

[Install]
WantedBy=multi-user.target


