# deployment/systemd/azl-sysproxy.service
[Unit]
Description=AZL Sysproxy (TCP bridge)
After=network.target
Wants=network.target
ConditionPathExists=/usr/local/bin/azl-sysproxy

[Service]
Type=simple
User=abdulrahman-alzalameh
WorkingDirectory=/home/abdulrahman-alzalameh/azl-language
EnvironmentFile=-/etc/default/azl
ExecStartPre=/bin/bash -lc 'mkdir -p .azl .azl/logs'
ExecStart=/bin/bash -lc 'SYSPROXY_TCP=${SYSPROXY_TCP:-127.0.0.1:9099} /usr/local/bin/azl-sysproxy 1>.azl/logs/sysproxy.out 2>.azl/logs/sysproxy.log'
Restart=always
RestartSec=1s
LimitNOFILE=8192
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectKernelTunables=yes
ProtectControlGroups=yes
ProtectClock=yes
LockPersonality=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
UMask=0077
TimeoutStopSec=10s
KillMode=process

[Install]
WantedBy=multi-user.target
