# deployment/systemd/azme-provider.service
[Unit]
Description=AZME Provider (AnythingLLM-compatible API)
After=network.target

[Service]
Type=simple
User=abdulrahman-alzalameh
WorkingDirectory=/home/abdulrahman-alzalameh/azl-language
Environment=AZME_PROVIDER_PORT=5000
# Backend selection for GPU models:
## Option A: Ollama (local GPU models)
Environment=AZME_BACKEND=ollama
Environment=AZME_OLLAMA_URL=http://127.0.0.1:11434
## Option B: OpenAI-compatible (e.g., vLLM/OpenAI)
#Environment=AZME_BACKEND=openai
#Environment=AZME_OPENAI_BASE_URL=http://127.0.0.1:8000
#Environment=AZME_OPENAI_API_KEY=
ExecStart=/usr/bin/env node scripts/azme_provider.js
Restart=always
RestartSec=2s
StandardOutput=append:/home/abdulrahman-alzalameh/azl-language/.azl/azme-provider.out
StandardError=append:/home/abdulrahman-alzalameh/azl-language/.azl/azme-provider.err
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
UMask=0077
LimitNOFILE=8192

[Install]
WantedBy=multi-user.target


